# src/.config/systemd/user/compton.service
[Unit]
Description=compton
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/compton
Restart=always

[Install]
WantedBy=xsession.target
