# src/.config/systemd/user/pnmixer.service
[Unit]
Description=pnmixer
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/pnmixer
Restart=always

[Install]
WantedBy=xsession.target
